# profile.d/bookmarks.bash
alias search-google='__yc_bookmark google'
alias search-google-image='__yc_bookmark google_image'
alias search-google-video='__yc_bookmark google_video'

alias search-ddg='__yc_bookmark duckduckgo'
alias search-ddg-image='__yc_bookmark duckduckgo_image'
alias search-ddg-video='__yc_bookmark duckduckgo_video'

alias search-youtube='__yc_bookmark youtube'
alias search-symbolhound='__yc_bookmark symbolhound'

alias search-mdn='__yc_bookmark mdn'
alias search-pydocs2='__yc_bookmark pydocs2'
alias search-pydocs3='__yc_bookmark pydocs3'
alias search-github='__yc_bookmark github'
alias search-bitbucket='__yc_bookmark bitbucket'
